# launchd/com.yasyf.claude-pool.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Reference template for the claude-pool user LaunchAgent.

  This file documents the plist that `clp service install` generates and writes
  to ~/Library/LaunchAgents/com.yasyf.claude-pool.plist. The installer fills in
  {{.Bin}}, {{.Log}} and {{.Path}} from the running binary; you normally do not
  edit this by hand. It is a *user* LaunchAgent (not a root LaunchDaemon)
  because credential refresh needs the user's login Keychain.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yasyf.claude-pool</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{.Bin}}</string>
        <string>daemon</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>{{.Log}}</string>
    <key>StandardErrorPath</key>
    <string>{{.Log}}</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{.Path}}</string>
    </dict>
</dict>
</plist>
